# supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:gsnova-paas]
command=/opt/gsnova/gsnova_paas_linux_amd64 -key $KEY -listen :8080
